# launchd/com.jake.autonomous-worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jake.autonomous-worker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/mikerodgers/Startup-Intelligence-OS/susan-team-architect/backend/scripts/jake_autonomous_worker.sh</string>
    </array>

    <!-- Run continuously — restart if it exits -->
    <key>KeepAlive</key>
    <true/>

    <!-- Start immediately on load -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Throttle rapid restarts (60s) -->
    <key>ThrottleInterval</key>
    <integer>60</integer>

    <!-- Stdout/stderr logs -->
    <key>StandardOutPath</key>
    <string>/Users/mikerodgers/.hermes/logs/autonomous_worker_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mikerodgers/.hermes/logs/autonomous_worker_stderr.log</string>

    <!-- Run as current user -->
    <key>UserName</key>
    <string>mikerodgers</string>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/mikerodgers/Startup-Intelligence-OS/susan-team-architect/backend</string>

    <!-- Nice priority — don't hog CPU -->
    <key>ProcessType</key>
    <string>Background</string>

    <!-- Environment — ensure PATH includes homebrew + node -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/mikerodgers</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
</dict>
</plist>
